# snapcast/rootfs/etc/services.d/snapcast/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Let's get started!"

## Run your program
exec /snapcast-script